# IT_infrastructre/ica0002_exam/roles/keepalived_exporter/templates/keepalived-exporter.service
[Unit]
Description=keepalived-exporter
Documentation=keepalived-exporter
After=network-online.target

[Service]
User=root
ExecStart=/usr/local/bin/keepalived-exporter

[Install]
WantedBy=multi-user.target
